# launchd/com.aiorch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.aiorch.plist — launchd agent for the AI Orchestrator watchdog daemon
  Installed to: ~/Library/LaunchAgents/com.aiorch.plist
  Generated by install.sh — do not edit directly; re-run install.sh to regenerate.
-->
<plist version="1.0">
<dict>

  <key>Label</key>
  <string>com.aiorch</string>

  <key>ProgramArguments</key>
  <array>
    <string>ORCH_HOME_PLACEHOLDER/bin/orch-daemon</string>
  </array>

  <!-- Start automatically when the user logs in -->
  <key>RunAtLoad</key>
  <true/>

  <!-- Keep alive: restart immediately if it exits -->
  <key>KeepAlive</key>
  <true/>

  <!-- Wait 10s before restarting after a crash -->
  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>ORCH_HOME_PLACEHOLDER/logs/daemon.log</string>

  <key>StandardErrorPath</key>
  <string>ORCH_HOME_PLACEHOLDER/logs/daemon-error.log</string>

  <!-- Ensure PATH includes Homebrew, user bins, and claude -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>HOME_PLACEHOLDER/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>HOME_PLACEHOLDER</string>
    <key>ORCH_HOME</key>
    <string>ORCH_HOME_PLACEHOLDER</string>
  </dict>

</dict>
</plist>
